# deploy/valet-v2.service
[Unit]
Description=V.A.L.E.T. kiosk
# The kiosk browser needs a graphical session, and the hub monitor needs the
# network to be up before the first scan.
After=network-online.target graphical.target
Wants=network-online.target

[Service]
Type=simple
User=pirc
Group=pirc
WorkingDirectory=/home/pirc/valet

# The kiosk reset launches Chromium, which needs access to the X display of the
# logged-in session. Adjust DISPLAY/XAUTHORITY to match your setup (these values
# are typical for a single auto-logged-in seat on :0).
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pirc/.Xauthority

# Bind the privileged ports 80 and 443 without running as root.
# CAP_NET_RAW/ADMIN are used when calling arp-scan
AmbientCapabilities=CAP_NET_BIND_SERVICE,CAP_NET_RAW,CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true

# -addr / -kiosk-url keep the kiosk browser on the loopback listener, so only
# its requests can trigger a kiosk reset. -http-addr / -https-addr expose the
# app to the LAN, where requests are non-loopback and cannot reset the kiosk.
# The self-signed TLS cert is generated under data/ on first start.
ExecStart=/home/pirc/valet/valet-v2 \
    -kiosk \
    -addr 127.0.0.1:3000 \
    -kiosk-url http://127.0.0.1:3000 \
    -http-addr :80 \
    -https-addr :443 \
    -data /home/pirc/valet/data/accounts.db \
    -tls-cert /home/pirc/valet/data/tls-cert.pem \
    -tls-key /home/pirc/valet/data/tls-key.pem

Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
