# deploy/valet-v2.service
[Unit]
Description=V.A.L.E.T. kiosk
# The kiosk browser needs a graphical session, and the hub monitor needs the
# network to be up before the first scan.
After=network-online.target graphical.target
Wants=network-online.target

[Service]
Type=simple
User=pirc
Group=pirc
WorkingDirectory=/home/pirc/valet

# The kiosk reset launches Chromium, which needs access to the X display of the
# logged-in session. Adjust DISPLAY/XAUTHORITY to match your setup (these values
# are typical for a single auto-logged-in seat on :0).
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pirc/.Xauthority

# Bind the privileged ports 80 and 443 without running as root.
# CAP_NET_RAW/ADMIN are used when calling arp-scan
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_NET_RAW CAP_NET_ADMIN
NoNewPrivileges=true

# The kiosk browser uses the HTTPS listener and receives an unguessable cookie
# at launch. That cookie lets it trigger resets without granting that ability
# to other clients on the same LAN listener.
# The self-signed TLS cert is generated under data/ on first start.
ExecStart=/home/pirc/valet/valet-v2 \
    -kiosk \
    -kiosk-url https://10.100.0.3 \
    -oauth-redirect-url https://10.100.0.3/login/complete \
    -data /home/pirc/valet/data/accounts.db \
    -tls-cert /home/pirc/valet/data/tls-cert.pem \
    -tls-key /home/pirc/valet/data/tls-key.pem

Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
